# contrib/inventory-agent.service
[Unit]
Description=inventory-agent service
Wants=network-online.target
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ocsinventory-agent --conf /etc/inventory-agent.conf

[Install]
WantedBy=multi-user.target
